# crontab
* * * * * yarn --cwd /js-api-client tape /packages/api-tests/tests/**.js -c creds/.test_config.json > /dev/stdout

