# crontab
*/15 * * * * yarn --cwd /js-api-client tape /js-api-client/packages/api-tests/tests/**.js -c /js-api-client/creds/.test_config.json > /dev/stdout

