# crontab
* * * * * yarn --cwd /js-api-clients tape packages/api-tests/tests/**.js > /dev/stdout
